# deploy/stock-agent.service
[Unit]
Description=AI Stock Trading Agent
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/stock-agent
EnvironmentFile=/home/ubuntu/stock-agent/.env
ExecStart=/usr/bin/java -Xmx2g -jar /home/ubuntu/stock-agent/target/stock-agent-1.0-SNAPSHOT.jar
Restart=always
RestartSec=15
StandardOutput=append:/home/ubuntu/stock-agent/logs/agent.log
StandardError=append:/home/ubuntu/stock-agent/logs/agent.log

[Install]
WantedBy=multi-user.target
